# etc/systemd/system/hardened-securonis-extra.service
[Unit]
Description=Apply Securonis Hardening Settings Once
After=multi-user.target
ConditionPathExists=!/var/lib/securonis-hardening-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/securonis-hardening-apply.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
